# srcs/FTPS/srcs/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid

[program:vsftpd]
command= /usr/sbin/vsftpd /etc/vsftpd/vsftpd.conf

[program:telegraf]
command=/usr/bin/telegraf -config /etc/telegraf.conf

[eventlistener:processes]
command=sh /stop-supervisor.sh
events=PROCESS_STATE_STOPPED, PROCESS_STATE_EXITED, PROCESS_STATE_FATAL